# etc/udev/rules.d/99-usb-midi.rules
ACTION=="add", SUBSYSTEM=="snd_seq", DRIVER=="snd_seq_midi", RUN+="/bin/sh /home/pi/figaro/bin/autoconnect.sh"